# wechatspider/shell/crontab
*/5 * * * * sh ~/wechat/script/start_comment.sh 2